# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create cct user if it doesn't exist
        if ! id -u cct >/dev/null 2>&1; then
            useradd --system --no-create-home --shell /bin/false cct
        fi

        # Create necessary directories
        mkdir -p /var/lib/cct/pocketbase
        mkdir -p /var/log/cct

        # Set ownership
        chown -R cct:cct /var/lib/cct
        chown -R cct:cct /var/log/cct

        # Set permissions
        chmod 750 /var/lib/cct/pocketbase
        chmod 750 /var/log/cct

        # Reload systemd
        systemctl daemon-reload

        # Enable and start PocketBase service
        systemctl enable cct-pocketbase.service
        systemctl start cct-pocketbase.service

        # Update desktop database
        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database /usr/share/applications
        fi

        echo ""
        echo "═══════════════════════════════════════════════════════════════"
        echo "  Claude Context Tracker (CCT) installed successfully!"
        echo "═══════════════════════════════════════════════════════════════"
        echo ""
        echo "✓ PocketBase is running at http://localhost:8090"
        echo "✓ Web interface is ready (no npm install needed!)"
        echo "✓ CLI tool available: cct"
        echo "✓ Desktop launcher added to Applications menu"
        echo ""
        echo "Get Started:"
        echo "  1. Open CCT from Applications menu → Development → Claude Context Tracker"
        echo "     OR run: cct-launcher"
        echo "     OR visit: http://localhost:8090"
        echo ""
        echo "  2. Create admin account at http://localhost:8090/_/"
        echo ""
        echo "  3. Use CLI commands:"
        echo "     cct status              # Show active project info"
        echo "     cct pull <project>      # Write context to CLAUDE.md"
        echo "     cct push <project> ...  # Save session summary"
        echo ""
        echo "Optional - Configure daemon to auto-track Claude Code:"
        echo "  sudo cp /usr/share/cct/cct-daemon.service.template /etc/systemd/system/cct-daemon.service"
        echo "  sudo nano /etc/systemd/system/cct-daemon.service  # Edit PROJECT_ID"
        echo "  sudo systemctl enable cct-daemon && sudo systemctl start cct-daemon"
        echo ""
        echo "Documentation: /usr/share/doc/cct/README.md"
        echo "═══════════════════════════════════════════════════════════════"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
